# kleinanzeigen-scanner/rootfs/etc/services.d/kleinanzeigen-scanner/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare TELEGRAM_API_TOKEN
declare KLEINANZEIGEN_API_USER
declare KLEINANZEIGEN_API_PW
declare PERSISTENCE_FILE
declare PERSISTENCE_FILE_BOT

## Get the 'message' key from the user config options.
export TELEGRAM_API_TOKEN=$(bashio::config 'TELEGRAM_API_TOKEN')
export KLEINANZEIGEN_API_USER=$(bashio::config 'KLEINANZEIGEN_API_USER')
export KLEINANZEIGEN_API_PW=$(bashio::config 'KLEINANZEIGEN_API_PW')
export PERSISTENCE_FILE=$(bashio::config 'PERSISTENCE_FILE')
export PERSISTENCE_FILE_BOT=$(bashio::config 'PERSISTENCE_FILE_BOT')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${TELEGRAM_API_TOKEN:="No API Token"}"

## Run your program
exec python3 /opt/kleinanzeigen-scanner/source/telegram_bot.py
